# daemon/zy4daemon.service
[Unit]
Description=Zy4Daemon Node Agent
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/zy4daemon
EnvironmentFile=/opt/zy4daemon/.env
ExecStart=/usr/local/bin/zy4daemon
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
